# schedule.cron
#57 19 * * * python3 /opt/scripts/tag_instance.py -ip -cp compartment-name -rg sa-saopaulo-1 -action add_defined -tag Periods.WeekDay="*, *, *, *, *, *, *, 1, *, *, *, *, *, *, *, *, *, *, 0, *, *, *, *, *"
#58 19 * * * python3 /opt/scripts/tag_instance.py -ip -cp compartment-name -rg sa-saopaulo-1 -action add_defined -tag Periods.Weekend="*, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *"
59 19 * * * python3 /opt/scripts/autorun.py -ip -rg sa-saopaulo-1 -tag Periods -a Down
#01 08 * * * python3 /opt/scripts/tag_instance.py -ip -cp compartment-name -rg sa-saopaulo-1 -action add_defined -tag Periods.WeekDay="*, *, *, *, *, *, *, 1, *, *, *, *, *, *, *, *, *, *, 0, *, *, *, *, *"
#02 08 * * * python3 /opt/scripts/tag_instance.py -ip -cp compartment-name -rg sa-saopaulo-1 -action add_defined -tag Periods.Weekend="*, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *, *"
03 08 * * * python3 /opt/scripts/autorun.py -ip -rg sa-saopaulo-1 -tag Periods -a Up
